# .config/systemd/user/greenclip.service
[Unit]
Description=Green clip

[Service]
ExecStart=/home/alen/.local/bin/greenclip daemon

[Install]
WantedBy=default.target
