# .config/systemd/user/redshift.service
[Unit]
Description=Music Player Daemon

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/redshift-gtk

[Install]
WantedBy=default.target
